# PocketVPN-Listener.service
[Unit]
Description=PocketVPN Pipe Listener

[Service]
ExecStart=/etc/PocketVPN/listen-pipe.sh

[Install]
WantedBy=multi-user.target
